# docker/nginx/crontab
5 10 * * * /usr/local/bin/php /var/www/html/web/wp/wp-cron.php >> /var/log/cron.log 2>&1 
